# systemd/owl-eye.target
[Unit]
Description=Owl-Eye Target
Wants=owl-eye.service

[Install]
WantedBy=default.target
